# jetson2_frying_ai/camera_autostart/sensing-camera.service
[Unit]
Description=SENSING Camera Driver Auto-loader
After=network.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/home/dkuyj/jetson-camera-monitor/jetson2_frying_ai/camera_autostart/camera_driver_autoload.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
