# jetson1_monitoring/jetson-monitor.service
[Unit]
Description=Jetson #1 Integrated Monitoring System
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=hr_dku_001
Group=hr_dku_001
WorkingDirectory=/home/hr_dku_001/jetson-food-ai/jetson1_monitoring
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/hr_dku_001/.Xauthority"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /home/hr_dku_001/jetson-food-ai/jetson1_monitoring/JETSON1_INTEGRATED.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security and resource limits
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=0

[Install]
WantedBy=graphical.target
