# jetson2_frying_ai/jetson2-ai.service
[Unit]
Description=Jetson #2 Frying AI Monitoring System
After=network.target sensing-camera.service
Wants=sensing-camera.service

[Service]
Type=simple
User=yjk
WorkingDirectory=/home/yjk/jetson-food-ai/jetson2_frying_ai
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/yjk/.Xauthority"
ExecStart=/usr/bin/python3 /home/yjk/jetson-food-ai/jetson2_frying_ai/JETSON2_INTEGRATED.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
